# agent/com.vaultbot.classifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vaultbot.classifier</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/apple/P1/Projects/Web/VaultBot/agent/start_worker.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/apple/P1/Projects/Web/VaultBot/agent</string>
    
    <key>StandardOutPath</key>
    <string>/Users/apple/P1/Projects/Web/VaultBot/agent/logs/worker.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/apple/P1/Projects/Web/VaultBot/agent/logs/worker.error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
